# main/modules/main_ofuscate/servicebackup_ofuscate/LOCK.service
[Unit]
Description=SERVICE_LOCK

[Service]
Type=simple
ExecStart=/LOCK/main/modules/servicebackup/LOCK.py

[Install]
WantedBy=multi-user.target