# main/modules/main_ofuscate/servicebackup_ofuscate/LOCK.service
[Unit]
Description=SERVICE_LOCK

[Service]
Type=simple
ExecStart=/home/joaog/LOCK/main/modules/servicemonitor/servicemonitor.py

[Install]
WantedBy=multi-user.target